# com.metaengine.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!--
    META ENGINE — macOS launchd Auto-Start
    =======================================
    This plist runs the Meta Engine at three times every trading day:
      - 9:21 AM ET  (pre-market brief — cached scans, zero-hour data)
      - 9:50 AM ET  (morning — 20 min of real market data)
      - 3:15 PM ET  (afternoon — power hour scan)
    
    Installation:
    1. Copy this file to ~/Library/LaunchAgents/
       cp com.metaengine.daily.plist ~/Library/LaunchAgents/
    
    2. Load the agent:
       launchctl load ~/Library/LaunchAgents/com.metaengine.daily.plist
    
    3. Verify it's loaded:
       launchctl list | grep metaengine
    
    4. To unload/stop:
       launchctl unload ~/Library/LaunchAgents/com.metaengine.daily.plist
    
    Note: The time is in LOCAL timezone. Adjust if you're not in ET.
    9:21 AM ET = Hour 9, Minute 21
    9:50 AM ET = Hour 9, Minute 50
    3:15 PM ET = Hour 15, Minute 15
    -->
    
    <key>Label</key>
    <string>com.metaengine.daily</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/chavala/Meta Engine/venv/bin/python3</string>
        <string>/Users/chavala/Meta Engine/run_meta_engine.py</string>
        <string>--force</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/chavala/Meta Engine</string>
    
    <key>StartCalendarInterval</key>
    <array>
        <!-- 9:21 AM ET — Pre-market brief -->
        <dict>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>21</integer>
        </dict>
        <!-- 9:50 AM ET — Morning scan (post-open) -->
        <dict>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>50</integer>
        </dict>
        <!-- 3:15 PM ET — Afternoon scan -->
        <dict>
            <key>Hour</key>
            <integer>15</integer>
            <key>Minute</key>
            <integer>15</integer>
        </dict>
    </array>
    
    <key>StandardOutPath</key>
    <string>/Users/chavala/Meta Engine/logs/launchd_stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/chavala/Meta Engine/logs/launchd_stderr.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>HOME</key>
        <string>/Users/chavala</string>
    </dict>
    
    <key>RunAtLoad</key>
    <false/>
    
    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
